# deploy/tencent-lighthouse/systemd/codewhale-runtime.service
[Unit]
Description=Mimofan Runtime API
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=mimofan
Group=mimofan
WorkingDirectory=/opt/whalebro
# Legacy /etc/deepseek is loaded first for old installs; /etc/mimofan wins.
EnvironmentFile=-/etc/deepseek/runtime.env
EnvironmentFile=-/etc/mimofan/runtime.env
ExecStart=/bin/sh -lc 'exec /home/mimofan/.cargo/bin/mimofan serve --http --host 127.0.0.1 --port "${MIMOFAN_RUNTIME_PORT:-${DEEPSEEK_RUNTIME_PORT:-7878}}" --workers "${MIMOFAN_RUNTIME_WORKERS:-${DEEPSEEK_RUNTIME_WORKERS:-2}}" --auth-token "${MIMOFAN_RUNTIME_TOKEN:-${DEEPSEEK_RUNTIME_TOKEN}}"'
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/mimofan/.mimofan /home/mimofan/.deepseek /opt/whalebro

[Install]
WantedBy=multi-user.target
