# deploy/tencent-lighthouse/systemd/codewhale-telegram-bridge.service
[Unit]
Description=Mimofan Telegram Phone Bridge
Wants=network-online.target mimofan-runtime.service
After=network-online.target mimofan-runtime.service

[Service]
Type=simple
User=mimofan
Group=mimofan
WorkingDirectory=/opt/mimofan/telegram-bridge
# Legacy /etc/deepseek is loaded first for old installs; /etc/mimofan wins.
EnvironmentFile=-/etc/deepseek/telegram-bridge.env
EnvironmentFile=-/etc/mimofan/telegram-bridge.env
ExecStart=/usr/bin/node /opt/mimofan/telegram-bridge/src/index.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/mimofan-telegram-bridge

[Install]
WantedBy=multi-user.target
