# deploy/tencent-lighthouse/systemd/codewhale-feishu-bridge.service
[Unit]
Description=Mimofan Feishu/Lark Phone Bridge
Wants=network-online.target mimofan-runtime.service
After=network-online.target mimofan-runtime.service

[Service]
Type=simple
User=mimofan
Group=mimofan
WorkingDirectory=/opt/mimofan/bridge
# Legacy /etc/deepseek is loaded first for old installs; /etc/mimofan wins.
EnvironmentFile=-/etc/deepseek/feishu-bridge.env
EnvironmentFile=-/etc/mimofan/feishu-bridge.env
ExecStart=/usr/bin/node /opt/mimofan/bridge/src/index.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/mimofan-feishu-bridge

[Install]
WantedBy=multi-user.target
